# Nginx_MySQL/cc_convt/data/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:converter]
command=/opt/RAIDA-CloudCoin/Coin_Converter_AES/converter
